# crontab
* * * * * python /cumple_bot.py >> /var/log/cron/cron.log 2>&1
# empty line
